# system-files/base/etc/profile.d/aliases.sh
#!/usr/bin/env bash

alias cat="bat --style=plain"